# 98-rafty-dd-one-from-udev.rules
KERNEL=="sr[0-9]*", ENV{DEVTYPE}=="disk", ENV{ID_CDROM_MEDIA_DVD}=="1", RUN+="/usr/sbin/systemctl start rafty-dd-dvd@dev-%k.service"
